# seriald/files/etc/uci-defaults/98-serial
cat > /etc/config/serial << EOF

config serial
	option baud   '115200'
	option device '/dev/ttyS0'
	option flow   'h'
	option socket ''
EOF

exit 0
