# i2cd-twists/files/etc/uci-defaults/50-hostname
waddr=$(uci get wireless.ap.macaddr)
hostname="flo-$(echo -n $waddr | tr -d ':' | cut -c 9-)"
common_hostname="flodevice"

domain_section=$(uci add dhcp domain)
uci batch << EOF
set system.@system[0].hostname='$hostname'
set network.wan.hostname='$hostname'
set dhcp.$domain_section.name='$common_hostname'
set dhcp.$domain_section.ip='192.168.3.1'
commit
EOF

exit 0
