# i2cd-twists/files/etc/uci-defaults/90-disble-dropbear-password-auth
uci batch << EOF
dropbear.@dropbear[0].PasswordAuth='off'
dropbear.@dropbear[0].RootPasswordAuth='off'
commit dropbear
EOF

exit 0
