# i2cd-twists/files/etc/uci-defaults/90-install-clashed-files
cat > /etc/banner << EOF
      :::::::::: :::        ::::::::
     :+:        :+:       :+:    :+:
    +:+        +:+       +:+    +:+
   :#::+::#   +#+       +#+    +:+
  +#+        +#+       +#+    +#+
 #+#        #+#       #+#    #+#
###        ########## ########
EOF

exit 0
