# i2cd-twists/files/etc/uci-defaults/98-support-external-opkg
echo "dest ext /tmp/vol/2" >> /etc/opkg.conf

exit 0
